# firstboot-debug.service
[Unit]
Description=First boot debug dump
After=multi-user.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/firstboot-debug.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target